# snapserver_pipe/rootfs/etc/services.d/snapserver/finish
#!/usr/bin/with-contenv sh
# Called when the service exits; keep simple for s6
exit 0
